# src/main/resources/puppet-repo.cron
#
# Puppet cron jobs.
#

# 
# Pull remote git repository every minute.
#
* * * * * root cd /etc/puppet-repo; git pull
